# redlooper.service
[Unit]
Description=RedLooper Application

[Service]
ExecStart=/usr/local/bin/redlooper

[Install]
WantedBy=default.target
